# services/usb_backup.timer
[Unit]
Description=USB backup timer

[Timer]
OnUnitActiveSec=600s
OnBootSec=10s

[Install]
WantedBy=timers.target
